# scripts/telegram-alert.service
[Unit]
Description=Telegram Alert Service for tv-webhook-trader
After=network-online.target tv-webhook-trader.service
Wants=network-online.target
PartOf=tv-webhook-trader.service

[Service]
Type=simple
User=ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/tv-webhook-trader/scripts/telegram_alert.py
Restart=always
RestartSec=10

# Environment - load from file so tokens stay out of the unit
EnvironmentFile=/home/ubuntu/tv-webhook-trader/.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-alert

[Install]
WantedBy=multi-user.target
